# subversion-stack/sources/crontab
@reboot svnserve -d -r /var/svn/repositories